# systemd/user/auto-sync-notes.timer
[Unit]
Description=Run auto note sync every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
