# btrbk/btrbk-backups.service
[Unit]
Description="Run btrbk"

[Service]
Type=oneshot
ExecStart=/usr/bin/btrbk run
